# alpine-nodejs/nodejs.ini
[supervisord]
nodaemon=true

[program:node]
command = /usr/bin/node app.js
directory = /var/lib/app/
user = root
autostart = true
stdout_logfile = /var/log/nodejs/out.log
stderr_logfile = /var/log/nodejs/err.log
